# mt5710-5g-connect.service
[Unit]
Description=Bring up TD Tech MT5710 5G NCM connection
After=systemd-udev-settle.service NetworkManager.service
Wants=NetworkManager.service
Conflicts=ModemManager.service

[Service]
Type=simple
Environment=MT5710_APN=ctnet
Environment=MT5710_CHECK_INTERVAL=20
Environment=MT5710_FAILURE_LIMIT=3
Environment=MT5710_RETRY_DELAY=15
ExecStart=/usr/local/sbin/mt5710-5g-connect.sh
Restart=always
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
